# inventory/services/fortress-service-demo.quadlet.d/web.container
[Container]
Notify=false

[Service]
RestartSec=5
